# timeclock/rootfs/etc/s6-overlay/s6-rc.d/postgres/run
#!/usr/bin/with-contenv bashio
# Bundled PostgreSQL 16 — loopback only. The cluster lives on /data (persisted
# across add-on updates); prepare.sh initdb's it on first boot.
export PGDATA=/data/postgres

bashio::log.info "Starting bundled PostgreSQL on 127.0.0.1:5432"
exec s6-setuidgid postgres \
  postgres -D "${PGDATA}" \
    -c listen_addresses=127.0.0.1 \
    -c unix_socket_directories=/run/postgresql \
    -c logging_collector=off
