# timeclock/rootfs/etc/s6-overlay/s6-rc.d/prepare/up
#!/command/execlineb -P
/etc/timeclock/prepare.sh
