# timeclock/rootfs/etc/s6-overlay/s6-rc.d/migrate/up
#!/command/execlineb -P
/etc/timeclock/migrate.sh
